# dockerfile/php/crontab
*/5 * * * * /usr/local/bin/php /var/www/sspanel/xcat Cron
